# build/http_handlers.service
[Unit]
Description=http_handlers

[Service]
Type=simple
User=root
ExecStart=/home/ec2-user/http_handlers

[Install]
WantedBy=multi-user.target
